# init.d/wdmd.service
[Unit]
Description=Watchdog Multiplexing Daemon
After=syslog.target

[Service]
Type=forking
ControlGroup=cpu:/
ExecStart=/lib/systemd/systemd-wdmd start
ExecStop=/lib/systemd/systemd-wdmd stop

[Install]
WantedBy=multi-user.target

